# .local/service/mailwatch/run
#!/bin/sh

#On startup, check for any new mails
mailsync

exec mailwatch cullenrss@gmail.com
