# .local/service/hlwm_idle_fullscreen/run
#!/bin/sh

exec ~/.config/herbstluftwm/idle_fullscreen.sh
