# .local/service/hlwm_idle_focus/run
#!/bin/sh

exec ~/.config/herbstluftwm/idle_focus.sh
