# .local/service/xautolock/run
#!/bin/sh

exec xautolock -locker 'sudo zzz -h'
