# .local/service/xautolock/run
#!/bin/sh

exec xautolock -time 5 -locker 'lock' -killtime 10 -killer 'sudo zzz -H' -notify 120 -notifier dim-screen.sh -detectsleep
