# .local/service/mpd-mpris/run
#!/bin/sh

exec mpd-mpris -no-instance -host localhost -port 6600
