# .local/service/barrier/run
#!/bin/sh

exec barrierc -f 192.168.0.8
