# .local/service/picom/run
#!/bin/sh

exec picom --backend glx --vsync --glx-no-stencil
