# .local/service/dbus-session/run
#!/bin/sh

exec dbus-daemon --session --nofork --nopidfile
