# .local/service/dbus-session/run
#!/bin/sh

exec dbus-daemon --session --nofork --nopidfile --address=unix:path=/run/usr/1000/bus
