# .local/service/hlwm_idle/run
#!/bin/sh

exec ~/.config/herbstluftwm/idle.sh
