# .local/service/xss-lock/run
#!/bin/sh

exec xss-lock -n dim-screen.sh -l -- lock
